# schedule-root.cron
# m h dom mon dow command
@reboot /usr/bin/motion -b
